# services/zero-thermometer.service
[Unit]
Description=Zero Thermometer Service
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/zero-thermometer
ExecStart=/home/admin/zero-thermometer/venv/bin/python /home/admin/zero-thermometer/main.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
